# root/etc/service/tlm_server/run
#!/bin/bash

# Start tlm_server
if ! pid=$(pgrep tlm_server);
then
  service start_tlmserver start
fi
